# templates/tomcat-init-centos.erb
#!/bin/bash
# description: Tomcat Start Stop Restart
# processname: tomcat
# chkconfig: 234 20 80
JAVA_HOME=/usr/lib/jvm/jre-1.7.0-openjdk.x86_64
export JAVA_HOME
PATH=$JAVA_HOME/bin:$PATH
export PATH
CATALINA_HOME="<%= @tomcat_home %>"

JAVA_OPTS="<%= @java_opts %>"

case $1 in
start)
# MC OOTB 20141010 - restarting doesn't work unless tomcat
# temp dirs are cleared down, copied this from clean_tomcat.sh
rm -rf $CATALINA_HOME/temp/Alfresco $CATALINA_HOME/work/Catalina/localhost/alfresco
rm -rf $CATALINA_HOME/work/Catalina/localhost/share
rm -rf $CATALINA_HOME/work/Catalina/localhost/awe
rm -rf $CATALINA_HOME/work/Catalina/localhost/wcmqs

su tomcat -c "sh $CATALINA_HOME/bin/startup.sh"
;; 
stop)   
su tomcat -c "sh $CATALINA_HOME/bin/shutdown.sh"
;; 
restart)
su tomcat -c "sh $CATALINA_HOME/bin/shutdown.sh"
su tomcat -c "sh $CATALINA_HOME/bin/startup.sh"
;; 
esac    
exit 0

