# templates/tomcat-init-centos.erb
#!/bin/bash
# description: Tomcat Start Stop Restart
# processname: tomcat
# chkconfig: 234 20 80
JAVA_HOME=/usr/lib/jvm/jre
export JAVA_HOME
PATH=$JAVA_HOME/bin:$PATH
export PATH
CATALINA_HOME="<%= @tomcat_home %>"

JAVA_OPTS="<%= @java_opts %>"
export JAVA_OPTS

function kill_tomcat {
	kill -9 `ps ax | grep tomcat | grep java | cut -c1-5`
}

case $1 in
start)
# MC OOTB 20141010 - restarting doesn't work unless tomcat
# temp dirs are cleared down, copied this from clean_tomcat.sh
rm -rf $CATALINA_HOME/temp/Alfresco $CATALINA_HOME/work/Catalina/localhost/alfresco
rm -rf $CATALINA_HOME/work/Catalina/localhost/share
rm -rf $CATALINA_HOME/work/Catalina/localhost/awe
rm -rf $CATALINA_HOME/work/Catalina/localhost/wcmqs

cd $CATALINA_HOME

su tomcat -c "sh $CATALINA_HOME/bin/startup.sh"
;; 
stop)   
su tomcat -c "sh $CATALINA_HOME/bin/shutdown.sh"
kill_tomcat
;; 
restart)
su tomcat -c "sh $CATALINA_HOME/bin/shutdown.sh"
kill_tomcat
su tomcat -c "sh $CATALINA_HOME/bin/startup.sh"
;; 
esac    
exit 0

